# backend/bittan/crontab
* * * * * cd ~/../bittan && /usr/local/bin/python manage.py run_cleaner >> bittancronoutput.txt 2>&1
